# ringlight-monitor.service
[Unit]
Description=RingLight Monitor - Automatic screen ring light
Documentation=https://github.com/yourusername/ringlight
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/ringlight-monitor-wrapper
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
