# systemd/ringlight-monitor.service
[Unit]
Description=RingLight Monitor - Auto-enable ring light for webcam/Howdy
Documentation=https://github.com/yourusername/ringlight
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/ringlight-monitor
Restart=on-failure
RestartSec=5

# Resource limits
MemoryMax=50M
CPUQuota=10%

[Install]
WantedBy=graphical-session.target
