# systemd/ringlight-monitor.service
[Unit]
Description=RingLight Monitor - Auto-enable ring light when webcam active
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ringlight-monitor -v
Restart=on-failure
RestartSec=3

# Can be overridden with systemctl --user edit ringlight-monitor
# Example overrides:
# ExecStart=/usr/local/bin/ringlight-monitor -v -s DP-1,HDMI-A-1 -a "-c FFFFFF -b 80"

[Install]
WantedBy=default.target
